# ePort/vending-machine.service
[Unit]
Description=Vending Machine Controller Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/soap
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 -m ePort.main
Restart=always
RestartSec=5

# Give time for hardware to initialize on boot
ExecStartPre=/bin/sleep 3

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
